# autorebooter.timer
[Unit]
Description=Run AutoRebooter periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target