# autorebooter.timer
[Unit]
Description=Run AutoRebooter periodically

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target